# config/config.symfolder/autostart/goa-replace.desktop
[Desktop Entry]
Name=GOA Replace
Exec=/usr/libexec/goa-daemon --replace
NoDisplay=true
Terminal=false
Type=Application

